# etc/profile.d/999-pulseaudio.sh
#!/bin/bash

# This scrips just starts the plus audio daemon for the user

pulseaudio -D
