# etc/profile.d/999-wine.sh
#!/bin/bash

# This script starts the wine prefix on login (except for root).

if [ "${USER}" != "root" ]; then
  wineboot
fi
